# crontab
PATH=/bin:/sbin:/usr/bin:/usr/local/bin

* * * * * cd /var/www/html; php -f cron.php  >> /var/log/cron.log 2>&1
